# autostart/bw.desktop
[Desktop Entry]
Name=bw update
Exec=bw.sh
Encoding=UTF-8
Type=Application
X-GNOME-Autostart-enabled=true
Icon=
